# release/ftp_upload.conf
[unix_http_server]
file=/var/run/supervisor.sock

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid 

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock 

; ftp上传服务
[program:ftp_upload]
directory=/home/toybrick/ftp_upload/release
command=/usr/bin/python3 main.py
autostart=true
user=toybrick
; nat
[program:nat]
directory=/home/toybrick/ftp_upload/release
command=/bin/bash iptables.sh
user=root
startsecs=0
autostart=true
autorestart=false
